# cron/cron.d/cron-jobs
* * * * * root docker exec gallery-wordpress-1 echo "Hello from app"
